# config/systemd/heybuddy.service
[Unit]
Description=heyBuddy AI Companion Service
Documentation=https://github.com/chrisperkles/heybuddy
After=network.target sound.target
Wants=network.target
RequiresMountsFor=/opt/heybuddy

[Service]
Type=simple
User=heybuddy
Group=heybuddy
WorkingDirectory=/opt/heybuddy
Environment=PYTHONPATH=/opt/heybuddy/src
Environment=ENVIRONMENT=production
Environment=SDL_VIDEODRIVER=dummy
Environment=SDL_AUDIODRIVER=alsa
EnvironmentFile=-/opt/heybuddy/.env

# Main service command
ExecStart=/opt/heybuddy/venv/bin/python /opt/heybuddy/src/main.py
ExecReload=/bin/kill -HUP $MAINPID

# Health check and status (optional - script may not exist)
# ExecStartPre=/opt/heybuddy/scripts/pre-start-check.sh
# WatchdogSec=60  # Only for Type=notify
# NotifyAccess=main  # Only for Type=notify

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryLimit=512M
CPUQuota=50%
TasksMax=50

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=no
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/heybuddy/data /opt/heybuddy/logs
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

# Audio access (required for PowerConf S330)
# Input access (required for button detection)
SupplementaryGroups=audio input
DevicePolicy=closed
DeviceAllow=/dev/snd/* rw
DeviceAllow=/dev/input/* rw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heybuddy

[Install]
WantedBy=multi-user.target